# hack/remote-debug/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:remote-debug-helper]
command=/usr/bin/remote-debug-helper
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sshd]
command=/usr/sbin/sshd -D -E /proc/1/fd/1
autostart=true
autorestart=true
